# fs_get_shell/etc/rc.local
#!/bin/bash

cleanup(){
    apt-get autoremove -y
    apt-get autoclean -y
    apt-get clean -y

    /sbin/shutdown -P -h now
}


trap cleanup INT TERM EXIT

apt-get update -y
apt-get dist-upgrade -y


/bin/bash



